# systemd/stream-server.service
[Unit]
Description=WebRTC Stream Server for Robot Dog
Documentation=https://github.com/yukebrillianth/ioh-robodog-webrtc
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ysc
WorkingDirectory=/home/ysc/ioh-robodog-webrtc
ExecStart=/home/ysc/ioh-robodog-webrtc/build/stream-server --config /home/ysc/ioh-robodog-webrtc/config.yaml
Environment=LD_LIBRARY_PATH=/home/ysc/ioh-robodog-webrtc/build/_deps/libdatachannel-build:/home/ysc/ioh-robodog-webrtc/build/_deps/spdlog-build
Environment=HOME=/home/ysc

# Restart policy
Restart=always
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=10

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stream-server

[Install]
WantedBy=multi-user.target
